# jamesmcguigan/modules/website_liatandco/templates/supervisor/liatandco.conf.erb
[program:liatandco.com]
environment=
command=/usr/bin/node /var/www/liatandco.com/Server.js NODE_ENV=production PORT_HTTP=3003
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true