# jamesmcguigan/modules/website_jamesmcguigan/templates/supervisor/jamesmcguigan.conf.erb
[program:jamesmcguigan.com]
environment=NODE_ENV=<%= @node_env %>
command=/usr/bin/node /var/www/jamesmcguigan.com/WebsiteCluster.js
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true