# jamesmcguigan/modules/website_jamesmcguigan_infographic/templates/supervisor/infographic.jamesmcguigan.conf.erb
[program:infographic.jamesmcguigan.com]
environment=NODE_ENV=<%= @node_env %>,PORT_HTTP=4000,PORT_HTTPS=4001
command=/usr/bin/node /var/www/infographic.jamesmcguigan.com/InfographicGenerator.js
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true